# systemd/radius-renew.service
[Unit]
Description=radius renew certs

[Service]
Type=oneshot
ExecStart=/etc/raddb/systemd/renew.sh
